# archlinux/roles/unbound/files/systemd-resolved.service.d/00-unit.conf
[Unit]
Requires=unbound.service consul@.service
After=unbound.service consul@.service

PartOf=unbound.service
