# archlinux/roles/fleet/files/fleet.service
[Unit]
Description=Fleetd Daemon

Requires=network-online.target

Wants=fleet.socket
After=fleet.socket


[Service]
#User=nobody
#Group=nobody

# applies user and group to ExecStart command only
PermissionsStartOnly=true

EnvironmentFile=/etc/environment

Environment=GOMAXPROCS=2

Restart=always

ExecStartPre=/usr/bin/mkdir -p /run/fleet
ExecStartPre=/usr/bin/chown root:root /run/fleet
ExecStartPre=/usr/bin/chmod 700 /run/fleet

ExecStart=/usr/bin/fleetd


[Install]
WantedBy=multi-user.target
