# archlinux/roles/systemd/files/system/ansible.timer
[Install]
WantedBy=multi-user.target

[Timer]
OnBootSec=10s
