# archlinux/roles/systemd/files/system/configure.timer
[Timer]
OnBootSec=0
WantedBy=multi-user.target
