# archlinux/roles/systemd/files/system/configure.timer
[Install]
WantedBy=multi-user.target

[Timer]
OnBootSec=0
