# archlinux/roles/systemd/files/system/environment.service
[Unit]
Description=Generate Environment files
After=systemd-resolved.service systemd-networkd-wait-online.service network-online.target
Requires=systemd-resolved.service systemd-networkd-wait-online.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
#Restart=on-failure
#RestartSec=5

ExecStart=/usr/local/bin/environment.py

[Install]
WantedBy=multi-user.target
