# archlinux/roles/systemd/files/system/systemd-resolved.service.d/00-consul.conf
[Unit]
After=consul.service
PartOf=consul.service
