# archlinux/roles/tmpwatch/files/tmpwatch.service
[Unit]
Description=tmpwatch

[Service]
ExecStart=/usr/bin/tmpwatch 168 /tmp
