# archlinux/roles/unbound/files/unbound.service.d/00-unit.conf
[Unit]
Requires=systemd-resolved.service consul@.service

PartOf=consul@.service

Before=systemd-resolved.service
After=consul@.service
