# archlinux/roles/unbound/files/unbound.service.d/00-unit.conf
[Unit]
Requires=consul@.service systemd-resolved.service
After=consul@.service
Before=systemd-resolved.service
PartOf=systemd-resolved.service
